# src/s6/loki/run
#!/command/execlineb -P

/stack/loki/loki --config.file=/stack/loki/config.yaml
