# src/s6/tempo/run
#!/command/execlineb -P

/stack/tempo/tempo --config.file=/stack/tempo/config.yaml
