# src/s6/opentelemetry-stack-exporter/run
#!/command/execlineb -P

python3 /stack/opentelemetry-stack-exporter/main.py
