# src/s6/otelcol-contrib/run
#!/command/execlineb -P

/stack/otelcol-contrib/otelcol-contrib --config=/stack/otelcol-contrib/config.yaml
