# src/s6/victoria-metrics/run
#!/command/execlineb -P

/stack/victoria-metrics/victoria-metrics-prod \
    -storageDataPath=/data/victoria-metrics \
    -retentionPeriod=90d \
    -usePromCompatibleNaming
