# src/s6/grafana/run
#!/command/execlineb -P

export GF_PATHS_HOME /data/grafana
export GF_PATHS_DATA /data/grafana/data
export GF_PATHS_PLUGINS /data/grafana/plugins

cd /stack/grafana
./bin/grafana server
