# src/s6/grafana/run
#!/command/execlineb -P

cd /stack/grafana
./bin/grafana server
