# src/s6/pyroscope/run
#!/command/execlineb -P

/stack/pyroscope/pyroscope --config.file=/stack/pyroscope/config.yaml
